# files/conf/consul.ini
[program:consul]
command=consul agent -config-dir=/etc/consul.d/server
stopsignal=INT
user=consul
stdout_logfile=/var/log/consul.log
autostart=true
autorestart=true
